# etc/rc.local
#!/bin/ash -e

/lib/rtl/init.sh &
/lib/aqr/init.sh &

exit 0
